# conf/supervisord_sample.conf
[unix_http_server]
file=/tmp/never_supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/path/to/log/files/never_supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/never_supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/never_supervisor.sock ; use a unix:// URL  for a unix socket


[program:never_server]
command = /path/to/virtualenvs/never-lib/bin/python never_server.py
directory = /path/to/project/never
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /path/to/logs/never.log
